# backend2/hancer-php-backend.service
[Unit]
Description=Hancer Law PHP Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/home/s/Belgeler/hancer.av.tr/dhavtr/backend2
ExecStart=/usr/bin/php -S 0.0.0.0:8000 index.php
Restart=always
RestartSec=10
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PHP_CLI_SERVER_WORKERS=4

[Install]
WantedBy=multi-user.target
